# deploy/seeduca.socket
[Unit]
Description=gunicorn blog socket

[Socket]
ListenStream=/run/chatbackend.socket

[Install]
WantedBy=sockets.target